# mDNSResponder/mDNSMacOSX/com.apple.mDNSResponderHelper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.apple.mDNSResponderHelper.reloaded</string>
  <key>OnDemand</key>
  <true/>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/sbin/mDNSResponderHelper</string>
  </array>
  <key>MachServices</key>
  <dict>
    <key>com.apple.mDNSResponder_Helper</key>
    <true/>
  </dict>
  <key>POSIXSpawnType</key>
  <string>Interactive</string>
</dict>
</plist>
